# python/src/smax-daemon/selector-smax-daemon.service
[Unit]
Description=Selector SMAX Daemon

[Service]
Type=notify
NotifyAccess=all
Environment=PYTHONUNBUFFERED=1
KillSignal=SIGINT
WorkingDirectory=/home/smauser/.config/systemd/user/selector-smax-daemon
ExecStart=/home/smauser/.config/systemd/user/selector-smax-daemon/on-start.sh

[Install]
WantedBy=default.target
